# sample.service.j2
[Unit]
Description=Simple service with API written in Go
After=network.target

[Service]
Type=simple
ExecStart=/home/vagrant/sample
PIDFile=/var/run/sample

[Install]
WantedBy=default.target
